# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
loglevel=info

[program:mysql]
command=/usr/bin/mysqld_safe
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/mysql.out.log
stderr_logfile=/var/log/supervisor/mysql.err.log

[program:backend]
command=python3 /app/back/main.py
autostart=true
autorestart=true
priority=30
environment=DB_HOST="127.0.0.1",DB_PORT="3306",DB_USER="recuser",DB_PASSWORD="recpass",DB_NAME="travel_rec"
stdout_logfile=/var/log/supervisor/backend.out.log
stderr_logfile=/var/log/supervisor/backend.err.log

[program:frontend]
command=streamlit run /app/front/app.py --server.port=8501 --server.address=0.0.0.0 --server.headless=true
autostart=true
autorestart=true
priority=40
environment=BACKEND_URL="http://127.0.0.1:8000"
stdout_logfile=/var/log/supervisor/frontend.out.log
stderr_logfile=/var/log/supervisor/frontend.err.log
